# systemd/node-health.timer
[Unit]
Description=Run node healthcheck every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
Unit=node-health.service

[Install]
WantedBy=timers.target
